# .flyio/crontab
# Immunity CRON tier - schedules to call the API tier on internal DNS.
# Loaded by supercronic in the immunity-cron Fly machine.

# Mock network ticker: every minute, simulates one round of activity.
* * * * * curl -sS -X POST -H "X-CRON-TOKEN: ${CRON_TOKEN}" http://immunity-api.internal/v1/internal/tick >/dev/null
